# t/recipes/checks/systemd/native/prerequisites/scripts-skip-systemd-native/build-spec/debian/scripts-skip-systemd-native-with.postinst
#!/bin/sh

set -e

#DEBHELPER#

invoke-rc.d
invoke-rc.d --skip-systemd-native

exit 0
