# t/recipes/checks/systemd/systemd-general/build-spec/debian/systemd-general.sourced.service
[Unit]
After=network.target
Documentation=man:sourced(1)

[Service]
ExecStart=/usr/bin/test

[Install]
WantedBy=sleep.target
