# testset/maintainer-scripts/debian/preinst
#!/bin/sh

update-rc.d foo remove
update-rc.d bar defaults
