# t/tests/systemd-general/debian/debian/test.service
[Unit]
After=network.target \
syslog.target

[Service]
ExecStartPre = /bin/true
ExecStart=/usr/bin/test

[Install]
WantedBy=multi-user.target
