# t/tags/checks/apache2/apache2-webapplications-general/debian/phpmyfoo.postrm
#!/bin/sh

set -e

if [ "$1" = 'configure' ] ; then
        a2disconf -q phpmyfoo
fi

invoke-rc.d apache2 reload

#DEBHELPER#

