# t/tests/init.d-general/debian/debian/init.d-general-bugs.postinst
#!/bin/sh

set -e

#DEBHELPER#

if [ "$1" = "configure" ] ; then
    # duplicate call as debhelper will add one for us
    update-rc.d init.d-general-bugs word
fi
