# t/tests/init.d-general/debian/debian/init.d-general-bugs.init
#!/bin/sh

set -e

### BEGIN INIT INFO
# Provides:          init.d-general-bugs
# Required-Start:
# Required-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: Do nothing at all
### END INIT INFO

# start stop force-reload restart status
