# t/recipes/checks/init-d/init.d-symlink/build-spec/debian/postinst
#!/bin/sh
set -e
update-rc.d lintian defaults >/dev/null
#DEBHELPER#
