# t/recipes/checks/systemd/systemd-complex-service-file/build-spec/debian/test3.service
[Unit]
After=
After = syslog.target

[Install]
WantedBy=multi-user.target
