# t/tags/checks/nmu/legacy-scripts/debian/scripts.conffiles
/etc/init.d/lsb-broken
/etc/init.d/lsb-other
/etc/init.d/no-lsb
/etc/X11/Xsession.d/xsession-test
/etc/csh/login.d/csh-foo
/etc/fish.d/fish-foo
