# t/recipes/checks/init.d/init.d-symlink/build-spec/debian/postrm
#!/bin/sh
set -e
if [ "$1" = purse ] ; then
    update-rc.d lintian remove >/dev/null
fi
#DEBHELPER#
