# t/tags/checks/systemd/systemd-service-file-pidfile-refers-to-var-run/debian/systemd-service-file-pidfile-refers-to-var-run-good.service
[Unit]
After=network.target
Documentation=https://example.com/

[Service]
ExecStart=/bin/test
PIDFile=/var/run/$NAME.pid

[Install]
WantedBy=multi-user.target
