# t/tests/apache2-webapplications-general/debian/debian/phpmyfoo.postinst
#! /bin/sh

set -e

if [ "$1" = 'configure' ] ; then
	a2enconf -q phpmyfoo
fi

invoke-rc.d apache2 reload

#DEBHELPER#
