# t/tests/legacy-scripts/debian/debian/postrm
#!/bin/sh -e

if [ "$1" = purge ] ; then
    update-rc.d lsb-broken remove >/dev/null
    update-rc.d no-lsb remove >/dev/null
    update-rc.d skeleton remove >/dev/null
    update-rc.d lsb-other remove >/dev/null
    update-rc.d lsb-other remove >/dev/null
fi
