# t/tests/scripts-missing-versioned-depends-on-init-system-helpers/debian/debian/scripts-missing-versioned-depends-on-init-system-helpers-without.postinst
#!/bin/sh

set -e

#DEBHELPER#

update-rc.d
update-rc.d defaults-disabled

exit 0
