# t/tests/init.d-lsb-headers/debian/debian/init.d-lsb-headers.init
#!/bin/sh

set -e

### BEGIN INIT INFO
# Provides:          nothing $bad_facility
# Required-Start:
# Required-Stop:
# Default-Start: 0 6 X
# Default-Stop: S Y
# Short-Description: Do nothing at all
### END INIT INFO

case "$1" in
    start|stop|force-reload|restart|status|*)
	echo hello world!
    ;;
esac
