# t/tests/scripts-calls-init-script/debian/debian/postrm
#!/bin/sh
set -e

if [ "remove" = "$1" ] ; then
    update-rc.d self-invoke remove >/dev/null
fi

#DEBHELPER#
