# t/tests/systemd-complex-service-file/debian/debian/test5.service
[Unit]
Documentation=https://www.example.com/

[Service]
Type=oneshot
ExecStart=/usr/bin/test
