# t/tests/systemd-complex-service-file/debian/debian/test2.service
.include /lib/systemd/system/test.service

[Unit]
After=
After=syslog.target
