# t/recipes/checks/init-d/maintainer-script/scripts-calls-init-script/build-spec/debian/postinst
#!/bin/sh
set -e

update-rc.d self-invoke defaults >/dev/null
/etc/init.d/self-invoke start

#DEBHELPER#
